# pandora/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare message

## Get the 'message' key from the user config options.
message=$(bashio::config 'message')

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "${message}"
bashio::log.info "Username: ${username}"

mkdir ~/.config
mkdir ~/.config/pianobar

echo "user = " $(bashio::config 'username') >> ~/.config/pianobar/config
echo "password = " $(bashio::config 'password') >> ~/.config/pianobar/config
echo "audio_quality = high" >> ~/.config/pianobar/config

contents=$(cat ~/.config/pianobar/config)

bashio::log.info "file contents: ${contents}"

## Run your program
exec pianobar
