# pandora/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare message

## Get the 'message' key from the user config options.
message=$(bashio::config 'message')
username=$(bashio::config 'username')
password=$(bashio::config 'password')
autostart_station=$(bashio::config 'autostart_station')

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "Building Pianobar config."

mkdir ~/.config
mkdir ~/.config/pianobar


# copy this into libao config -- test if needed
cat >>~/.libao <<EOL
default_driver=default
EOL

#copy this into alsa config -- not needed
# cat >>~/.asoundrc <<EOL
# # Use pulseaudio by default
# pcm.!default {
#   type pulse
#   fallback "sysdefault"
#   hint {
#     show on
#     description "Default ALSA Output (currently PulseAudio Sound Server)"
#   }
# }

# ctl.!default {
#   type pulse
#   fallback "sysdefault"
# } 
# EOL



echo "user = ${username}" >> ~/.config/pianobar/config
echo "password = ${password}" >> ~/.config/pianobar/config
echo "audio_quality = high" >> ~/.config/pianobar/config
echo "autostart_station = ${autostart_station}" >> ~/.config/pianobar/config

configcontents=$(cat ~/.config/pianobar/config)
libaocontents=$(cat ~/.libao)
asoundcontents=$(cat /etc/asound.conf)

bashio::log.info "libao config file contents: ${libaocontents}"
bashio::log.info "config file contents: ${configcontents}"
bashio::log.info "/etc/asound.conf contents: ${asoundcontents}"

bashio::log.info "Alsa devices:"
aplay -L

bashio::log.info "Alsa cards:"
aplay -l

## Run your program

speaker-test -t wav -c 2 --nloops 1

exec pianobar
